# templates/uwsgi.conf.erb
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [<%= runlevel_start %>]
stop on runlevel [<%= runlevel_stop %>]

exec uwsgi --master --die-on-term --emperor <%= name %> --logto /var/log/uwsgi/<%= service %>.log